# packages/oscam/unit/oscam.service
[Unit]
Description=Oscam Services

[Service]
ExecStart=/usr/bin/oscam -d 128
Restart=always

[Install]
WantedBy=multi-user.target
