# packages/vdr/unit/VDR.target
[Unit]
Description=VDR Video Disc Recorder
Requires=graphical.target
Conflicts=rescue.service rescue.target
After=graphical.target
AllowIsolate=yes
